# udev/rules.d/70_wifi_powersave.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlp*", RUN+="/usr/bin/iw dev %k set power_save on"
